# supervisor-production.conf
; =============================================================================
; QuizSnap — Supervisor (PRODUCTION ONLY)
; =============================================================================
; SERVER: 48 GB RAM | 12 CPU — QuizSnap CPU budget: 8 cores
;
; Process budget (conservative):
;   Reverb:          1 process  (~0.5 CPU, ~300 MB RAM)
;   Queue default:   4 workers  (~1.5 CPU, ~600 MB RAM)
;   Queue imports:   1 worker   (~0.5 CPU, ~150 MB RAM)
;   PHP-FPM:         separate pool (48 max children)
;   Total QuizSnap artisan: 6 processes — leaves CPU for PHP-FPM under load
;
; Do NOT increase numprocs without running production-audit.sh during peak exam.
;
; INSTALL: /etc/supervisor/conf.d/quizsnap.conf
; Development: manual queue:work or database queue — not required locally.
; =============================================================================

[program:quizsnap-reverb]
process_name=%(program_name)s
command=php /var/www/quizsnap/artisan reverb:start --host=0.0.0.0 --port=8080
directory=/var/www/quizsnap
autostart=true
autorestart=true
user=www-data
numprocs=1
priority=10
redirect_stderr=true
stdout_logfile=/var/log/supervisor/quizsnap-reverb.log
stopwaitsecs=10
environment=HOME="/var/www/quizsnap",APP_ENV="production"

[program:quizsnap-queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/quizsnap/artisan queue:work redis --queue=default --sleep=2 --tries=3 --timeout=90 --max-time=3600 --memory=128
directory=/var/www/quizsnap
autostart=true
autorestart=true
user=www-data
numprocs=4
priority=20
redirect_stderr=true
stdout_logfile=/var/log/supervisor/quizsnap-queue-default.log
stopwaitsecs=3600
environment=HOME="/var/www/quizsnap",APP_ENV="production"

[program:quizsnap-queue-imports]
process_name=%(program_name)s
command=php /var/www/quizsnap/artisan queue:work redis --queue=imports,notifications --sleep=5 --tries=3 --timeout=300 --max-time=3600 --memory=192
directory=/var/www/quizsnap
autostart=true
autorestart=true
user=www-data
numprocs=1
priority=30
redirect_stderr=true
stdout_logfile=/var/log/supervisor/quizsnap-queue-imports.log
stopwaitsecs=3600
environment=HOME="/var/www/quizsnap",APP_ENV="production"

; Scheduler — use system cron (one line), not a busy-loop process:
; * * * * * www-data cd /var/www/quizsnap && php artisan schedule:run >> /dev/null 2>&1
