# scripts/vps/supervisor-quizsnap.conf
; QuizSnap VPS — copy to /etc/supervisor/conf.d/quizsnap.conf
; Adjust paths and user. Then: supervisorctl reread && supervisorctl update

[program:quizsnap-reverb]
process_name=%(program_name)s
command=php /var/www/quizsnap/artisan reverb:start --host=0.0.0.0 --port=8080
directory=/var/www/quizsnap
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/quizsnap-reverb.log
stopwaitsecs=10

[program:quizsnap-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/quizsnap/artisan queue:work redis --sleep=1 --tries=3 --max-time=3600 --timeout=120
directory=/var/www/quizsnap
autostart=true
autorestart=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/quizsnap-queue.log
stopwaitsecs=3600

; Scale queue workers: numprocs=4 on 6-core VPS is a good start for ~3000 quiz takers.
; Increase Reverb instances only with REVERB_SCALING_ENABLED=true + Redis.
