# scripts/vps/supervisor-quizsnap.conf
; QuizSnap VPS — copy to /etc/supervisor/conf.d/quizsnap.conf
; Default app path matches deploy.sh (APP_DIR=/srv/apps/quizsnap).
; Then: supervisorctl reread && supervisorctl update && supervisorctl restart quizsnap-reverb quizsnap-worker

[program:quizsnap-reverb]
process_name=%(program_name)s
command=php /srv/apps/quizsnap/artisan reverb:start --host=127.0.0.1 --port=8080
directory=/srv/apps/quizsnap
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/quizsnap-reverb.log
stopwaitsecs=10
environment=HOME="/srv/apps/quizsnap",APP_ENV="production"

[program:quizsnap-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /srv/apps/quizsnap/artisan queue:work redis --sleep=1 --tries=3 --max-time=3600 --timeout=120
directory=/srv/apps/quizsnap
autostart=true
autorestart=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/quizsnap-worker.log
stopwaitsecs=3600
environment=HOME="/srv/apps/quizsnap",APP_ENV="production"

; Cron (every minute):
; * * * * * www-data cd /srv/apps/quizsnap && php artisan schedule:run >> /dev/null 2>&1
